# ansible/roles/herlighet/templates/herlighet.service.j2
[Unit]
Description=herlighet

[Service]
Restart=always
User=herlighet
ExecStart=/opt/herlighet/bin/herlighet

[Install]
WantedBy=multi-user.target
